# deb-build/ovs-plugin/lib/systemd/system/docker-ovs-plugin.service
[Unit]
Description=Docker OVS Plugin
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/docker-ovs-plugin
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

